# 新建文件夹/excel_upload_system/excel_upload_system.conf
[program:excel_upload_system]
# 应用程序命令
command=/opt/excel_upload_system/venv/bin/gunicorn -c /opt/excel_upload_system/gunicorn_config.py main:app

# 应用程序工作目录
directory=/opt/excel_upload_system

# 启动用户（生产环境建议使用非root用户）
user=www-data

# 自动启动和重启配置
autostart=true            ; 当supervisor启动时，自动启动该应用
autorestart=true          ; 应用崩溃时自动重启
startsecs=5               ; 启动后等待5秒认为启动成功
startretries=3            ; 启动失败后重试次数
restart_delay=5           ; 重启间隔时间（秒）
restart_max_delay=60      ; 最大重启间隔时间

# 停止信号配置
stopsignal=TERM           ; 使用TERM信号停止进程
stopwaitsecs=300          ; 等待进程停止的最大秒数

# 进程数量控制
numprocs=1                ; 启动的进程数量

# 日志配置
stdout_logfile=/opt/excel_upload_system/logs/supervisor_stdout.log  ; 标准输出日志
stderr_logfile=/opt/excel_upload_system/logs/supervisor_stderr.log  ; 错误输出日志
stdout_logfile_maxbytes=50MB     ; 单个日志文件最大大小
stderr_logfile_maxbytes=50MB     ; 单个错误日志文件最大大小
stdout_logfile_backups=10        ; 标准日志备份数量
stderr_logfile_backups=10        ; 错误日志备份数量
loglevel=info                    ; 日志级别

# 环境变量设置
environment=PATH="/opt/excel_upload_system/venv/bin:%(ENV_PATH)s", \
            FLASK_ENV="production", \
            FLASK_APP="main.py", \
            PYTHONPATH="/opt/excel_upload_system"

# 进程优先级
priority=999

# 停止顺序组
stopasgroup=false         ; 停止时是否同时停止子进程组
killasgroup=false         ; 终止时是否强制终止子进程组

# 状态检查配置
redirect_stderr=true      ; 将stderr重定向到stdout

# 事件监听
events=true               ; 允许事件监听

# 进程启动超时
timeout=300               ; 进程启动超时时间（秒）

# 依赖项
# 如果有其他依赖服务，可以在这里配置
dependencies=            ; 依赖的其他程序名称